# scripts/com.edith.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.edith.ai</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>--import</string>
    <string>tsx/esm</string>
    <string>__INSTALL_DIR__/src/main.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__INSTALL_DIR__</string>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>
  <key>StandardOutPath</key>
  <string>__INSTALL_DIR__/.edith/logs/edith.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>__INSTALL_DIR__/.edith/logs/edith.stderr.log</string>
</dict>
</plist>
